# birdorsquirrel-supervisor.conf
[program:birdorsquirrel]
directory=/home/pi/BirdOrSquirrel/
command=python birdorsquirrel.py
autostart=true
autorestart=true
