# supervisor.conf
[program:birdorsquirrel]
command = python /app/birdorsquirrel.py
